# rootfs/etc/services.d/invoices-push/run
#!/usr/bin/execlineb -P

with-contenv
php /var/www/bin/console app:listener